# x250_system_config/udev_rules.d/10-local.rules
KERNEL=="card0", ACTION=="change", ENV{XAUTHORITY}="/home/patrick/.Xauthority", ENV{DISPLAY}=":0.0", RUN+="/home/patrick/scripts/dp-script.sh"
